# rust-proxy/rust-proxy.service
[Unit]
Description=Rust Reverse Proxy
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/homeroute/rust-proxy/target/release/rust-proxy
ExecReload=/bin/kill -HUP $MAINPID
Environment=PROXY_CONFIG_PATH=/var/lib/server-dashboard/rust-proxy-config.json
Environment=RUST_LOG=rust_proxy=info,tower_http=debug
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rust-proxy

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/server-dashboard
ReadOnlyPaths=/opt/homeroute/rust-proxy
ProtectHome=true

# Allow binding to privileged ports if needed
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
