# systemd/hr-apps.slice
[Unit]
Description=HomeRoute user applications slice
Before=slices.target

[Slice]
MemoryAccounting=yes
CPUAccounting=yes
TasksAccounting=yes
