# systemd/hr-edge.service
[Unit]
Description=HomeRoute Edge (HTTPS Proxy + TLS + ACME + Auth + Tunnel)
After=network-online.target hr-netcore.service
Wants=network-online.target hr-netcore.service

[Service]
Type=simple
ExecStart=/opt/homeroute/crates/target/release/hr-edge
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/opt/homeroute
EnvironmentFile=/opt/homeroute/.env
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hr-edge

# Ports 443 and 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/server-dashboard /opt/homeroute/data /run /tmp
ReadOnlyPaths=/opt/homeroute/crates/target /opt/homeroute/web/dist
ProtectHome=false

[Install]
WantedBy=multi-user.target
