# systemd/hr-netcore.service
[Unit]
Description=HomeRoute Network Core (DNS/DHCP/Adblock/IPv6)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/homeroute/crates/target/release/hr-netcore
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/opt/homeroute
EnvironmentFile=/opt/homeroute/.env
Restart=always
RestartSec=200ms
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hr-netcore

# Capabilities: bind port 53, raw sockets (DHCP/RA), network admin (nftables), broadcast (DHCP)
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_NET_ADMIN CAP_NET_BROADCAST

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/server-dashboard /opt/homeroute/data /var/log/rust-dns-dhcp /run /tmp
ReadOnlyPaths=/opt/homeroute/crates/target
ProtectHome=yes

[Install]
WantedBy=multi-user.target
